# supervisord.simple.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:landing]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/landing.err.log
stdout_logfile=/var/log/supervisor/landing.out.log
environment=PORT="8080",NODE_ENV="production"

[program:crawl4ai]
command=uvicorn server:app --host 0.0.0.0 --port 11235
directory=/app/crawl4ai-source/deploy/docker
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/crawl4ai.err.log
stdout_logfile=/var/log/supervisor/crawl4ai.out.log
environment=PYTHONUNBUFFERED="1"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log 