# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Landing page service
[program:landing]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/landing.err.log
stdout_logfile=/var/log/supervisor/landing.out.log
environment=PORT="8080",NODE_ENV="production"

# n8n service
[program:n8n]
command=n8n start
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/n8n.err.log
stdout_logfile=/var/log/supervisor/n8n.out.log
environment=N8N_HOST="0.0.0.0",N8N_PORT="5678",NODE_ENV="production"

# Crawl4AI service
[program:crawl4ai]
command=uvicorn server:app --host 0.0.0.0 --port 11235
directory=/app/crawl4ai-source/deploy/docker
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/crawl4ai.err.log
stdout_logfile=/var/log/supervisor/crawl4ai.out.log
environment=PYTHONUNBUFFERED="1",DISPLAY=":99"

# Nginx reverse proxy
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log 